# systemd/docker-services-disk-check.service
[Unit]
Description=Check docker-services host disk usage and notify on threshold
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=lightweight
WorkingDirectory=/home/lightweight/docker-services
ExecStart=/home/lightweight/docker-services/scripts/disk-check.sh
