# systemd/docker-services-startup-check.service
[Unit]
Description=Run docker-services startup validation and health summary
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=lightweight
WorkingDirectory=/home/lightweight/docker-services
Environment=MACHINE_ID=server01
ExecStart=/home/lightweight/docker-services/scripts/startup-check.sh
